# deployment/systemd/binance-extreme-funding-dashboard.service
[Unit]
Description=Binance extreme-funding dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/arbitrage
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/arbitrage/.venv/bin/python -m binance_extreme_funding.run_dashboard --host 127.0.0.1 --port 8770
Restart=always
RestartSec=15
Nice=10

[Install]
WantedBy=multi-user.target
